# files/system/etc/profile.d/lessopts.sh
# Set less default arguments
export LESS=-F # don't go into pager mode if there is less than one page anyway