# files/system/etc/profile.d/missing-env.sh
# Help yaml-streamly find libyaml
export PKG_CONFIG_LIBDIR=/home/linuxbrew/.linuxbrew/Cellar/libyaml/*/lib/pkgconfig